# net-analyzer/xython/files/xython-celery.init
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/python3"
command_args="-m celery -A xython worker"
PIDFILE="/run/xython-celery.pid"
start_stop_daemon_args="--background --make-pidfile --user xython:xython --pidfile ${PIDFILE}"

depend() {
	need localmount net
	use logger
}
